# scripts/backup-verify.timer
[Unit]
Description=Weekly backup verification

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
