# scripts/com.opensintunnel.healthcheck.timer
# SPDX-License-Identifier: MIT
[Unit]
Description=run OpenSIN tunnel health check every 30 seconds
Documentation=https://github.com/OpenSIN-AI/OpenSIN-Chat

[Timer]
OnBootSec=1min
OnUnitActiveSec=30s
AccuracySec=1s
Persistent=true
Unit=com.opensintunnel.healthcheck.service

[Install]
WantedBy=timers.target
