# ci/ci-webhook.service
[Unit]
Description=OpenSIN-Chat CI Webhook Receiver
After=network.target docker.service
Wants=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/OpenSIN-Chat
EnvironmentFile=/home/ubuntu/ci-webhook.env
ExecStart=/usr/bin/node /home/ubuntu/OpenSIN-Chat/ci/webhook-server.cjs
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
